# .config/systemd/user/systemd-list-updater.service
[Unit]
Description=Actualiza la lista de unidades habilitadas de systemd

[Service]
Type=oneshot
ExecStart=%h/.local/scripts/update-systemd-list.sh
