# .config/systemd/user/systemd-list-updater.path
[Unit]
Description=Monitorea cambios en los servicios habilitados

[Path]
PathChanged=/etc/systemd/system

[Install]
WantedBy=default.target
